# supervisor/project-api.ini
[program:project_api]
command=/home/server-name/.pyenv/versions/project-api/bin/gunicorn -w 1 -b 0.0.0.0:5000 --pythonpath='/home/server-name/.pyenv/versions/project-api/bin/python3.6,/home/server-name/.pyenv/versions/project-api/lib/python3.6/site-packages' project_api.wsgi:app
directory=/home/server-name/python/project-api
user=root
startsecs=3
stopwaitsecs=1
autostart=true
autorestart=true
stdout_logfile=/etc/supervisord/logs/project-api.log
stderr_logfile=/etc/supervisord/logs/project-api-err.log
environment=STAGING_CONFIG=/home/server-name/python/project-api/project_api/instance/staging.py